# addon/rootfs/etc/services.d/telegram-face-bot/run
#!/usr/bin/with-contenv bash
set -euo pipefail

# Read options from Home Assistant (support new and legacy key)
TOKEN=$(jq -r '(.TELEGRAM_BOT_API_TOKEN // .bot_token // empty)' /data/options.json)

if [ -z "${TOKEN:-}" ]; then
  echo "[ERROR] No token provided. Set 'TELEGRAM_BOT_API_TOKEN' in add-on options." >&2
  exit 1
fi

export TELEGRAM_BOT_API_TOKEN="$TOKEN"

cd /opt/telegram-face-bot
exec python3 /opt/telegram-face-bot/main.py


